# Docker-trac/compose/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:sites-available1]
command=/bin/bash -c "a2ensite trac.conf"
priority=1
autostart=true
startsecs=0
exitcodes=0
autorestart=false

[program:sites-available2]
command=/bin/bash -c "a2ensite svn.conf"
priority=2
autostart=true
startsecs=0
exitcodes=0
autorestart=false

[program:sites-available3]
command=/bin/bash -c "a2ensite git.conf"
priority=3
autostart=true
startsecs=0
exitcodes=0
autorestart=false

[program:sites-available4]
command=/bin/bash -c "service apache2 reload"
priority=4
autostart=true
startsecs=0
exitcodes=0
autorestart=false

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
priority=10
autorestart=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autorestart=true

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
autorestart=true

#[program:tracd]
#command=/usr/local/bin/tracd --port 8000 --basic-auth='*,/var/trac/conf/.htpasswd,trac520'  /var/trac
#autorestart=true
